# meta-celestica/meta-as58xx-cl/recipes-plat/ipmid/files/ipmid.service
[Unit]
Description=System IPMID Service

[Service]
ExecStart=/usr/local/bin/ipmid
Restart=on-failure

[Install]
WantedBy=multi-user.target

