# meta-celestica/meta-as58xx-cl/recipes-plat/plat-utils/files/power_monitor.service
[Unit]
Description=System Power Monitor Service
After=platform.service

[Service]
ExecStart=/usr/local/bin/power_monitor.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target

