# meta-celestica/meta-as58xx-cl/recipes-utils/openbmc-utils/files/bmc_wdt.service
[Unit]
Description=BMC Watchdog Service
After=systemd-logind.service

[Service]
Type=forking
ExecStart=/sbin/watchdog -t 30 /dev/watchdog
Restart=on-failure

[Install]
WantedBy=multi-user.target

