# meta-celestica/meta-as58xx-cl/recipes-core/fan-ctrl/fan-ctrl/fand.service
[Unit]
Description=Fan Control Service
After=power-on.service

[Service]
Type=forking
ExecStart=/usr/local/bin/setup-fan.sh start
Restart=on-failure

[Install]
WantedBy=multi-user.target

