# meta-celestica/meta-as58xx-cl/recipes-plat/ipmbd/files/ipmbd.service
[Unit]
Description=System IPMBD Service

[Service]
ExecStart=/usr/local/bin/ipmbd 1 6
Restart=on-failure

[Install]
WantedBy=multi-user.target

