# meta-celestica/meta-as58xx-cl/recipes-plat/openbmc-gpio/files/openbmc-gpio-1/openbmc_gpio.service
[Unit]
Description=OpenBmc GPIO Service

[Service]
ExecStart=/usr/local/bin/openbmc_gpio_setup.py
Restart=on-failure

[Install]
WantedBy=multi-user.target

